# systemd/quote.timer
[Unit]
Description=Timer for ink display quote on the :30 mark

[Timer]
OnCalender=*:30
Persisent=true

[Install]
WantedBy=timers.target
