# systemd/user/ink_weather.timer
[Unit]
Description=Timer for ink display quote on the :00 mark

[Timer]
OnCalendar=*:00
Persistent=true

[Install]
WantedBy=timers.target
