# systemd/user/ink_fact.timer
[Unit]
Description=Timer for ink display random fact on the :15 mark

[Timer]
OnCalendar=*:15
Persistent=true

[Install]
WantedBy=timers.target
