# systemd/user/ink_quote.timer
[Unit]
Description=Timer for ink display quote on the :30 mark

[Timer]
OnCalendar=*:30
Persistent=false

[Install]
WantedBy=timers.target
